# serv_er/supervisord.ini
[supervisord]
nodaemon=true

[program:flask]
command=uwsgi --ini /etc/uwsgi/uwsgi.ini
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/uwsgi.log
stderr_logfile=/var/log/uwsgi_err.log

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx.log
stderr_logfile=/var/log/nginx_err.log
